# systemd/discord-codex-bridge.service
[Unit]
Description=Discord Codex Backup Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/yang0/projectHome/disocord_codex
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/yang0/projectHome/disocord_codex/.venv/bin/python -m discord_codex_bridge --env-file /home/yang0/projectHome/disocord_codex/.env
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
