# systemd/discord-codex-bridge.service
[Unit]
Description=Discord Codex Backup Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/projectHome/discord-codex-bridge
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=%h/projectHome/discord-codex-bridge/src
ExecStart=%h/projectHome/discord-codex-bridge/.venv/bin/python -m discord_codex_bridge --env-file %h/projectHome/discord-codex-bridge/.env
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
